# apps/mail/mail/rootfs/etc/services.d/syslog/run
#!/usr/bin/env bash

echo "s6-supervise: spawning syslog process"

#exec syslogd -n -t -L -R fail2ban
exec syslogd -n -t
