# apps/mail/unbound/rootfs/etc/services.d/unbound/run
#!/usr/bin/env sh

echo "s6-supervise: spawning unbound process"

# Need to specify config file
exec unbound &> /dev/stdout
