# apps/antivirus/rootfs/etc/services.d/clamav-milter/run
#!/usr/bin/env sh

echo "[INFO] Starting clamav-milter"

# To test configuration:
#exec clamav-milter

exec clamav-milter &>/dev/null

exit 0
