# apps/mail/opendkim/rootfs/etc/services.d/opendkim/run
#!/usr/bin/env sh

echo "s6-supervise: spawning opendkim process"

# More verbose
#exec opendkim -f -v

exec opendkim -f
