# apps/mail/mail/rootfs/etc/services.d/dovecot/run
#!/usr/bin/env bash

logger "s6-supervise: spawning process"

exec dovecot -F
